# mqtt-bridge.service
[Unit]
Description=MQTT to PostgreSQL Bridge Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/MQTT
EnvironmentFile=/root/MQTT/.env
ExecStart=/usr/bin/python3 /root/MQTT/mqtt_to_postgres.py
Restart=always
RestartSec=10
StandardOutput=append:/root/MQTT/log/mqtt_bridge.log
StandardError=append:/root/MQTT/log/mqtt_bridge_error.log

[Install]
WantedBy=multi-user.target
